# .config/systemd/user/gammastep-refresh.timer
[Unit]
Description=Run gammastep refresh every hour

[Timer]
OnBootSec=1min
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target

